# deploy/supervisor-nexa-queue.conf
; Supervisor program config untuk queue worker NEXA production.
; Menggantikan `php artisan queue:listen` (dipakai composer script `dev`,
; cuma cocok untuk development) — production butuh proses daemon yang
; auto-restart kalau crash/di-reload deploy, itu tugas Supervisor di sini.
;
; Install:
;   sudo cp deploy/supervisor-nexa-queue.conf /etc/supervisor/conf.d/nexa-queue.conf
;   # sesuaikan `directory`, `user`, dan path `php` di bawah dulu
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start nexa-queue:*
;
; Restart queue worker setiap kali deploy kode baru (worker memegang kode
; lama di memory sampai di-restart):
;   sudo supervisorctl restart nexa-queue:*

[program:nexa-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/nexa/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/nexa
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/nexa/storage/logs/queue-worker.log
stopwaitsecs=3600
